# py/ops/itests/testdata/bundle1/simple.service
[Unit]
Description=simple

[Service]
ExecStart=/bin/sleep 1d
