# shipyard/shipyard/py/garage/examples/echod/echod.service
[Unit]
Description=echod

[Service]
Slice=machine.slice
ExecStart=/usr/bin/rkt run --pod-manifest ${POD_MANIFEST}
KillMode=mixed
Restart=always
