# py/ops/itests/testdata/bundle3/volume.service
[Unit]
Description=volume

[Service]
ExecStart=/bin/sleep 1d
