# py/ops/itests/testdata/bundle2/replicated.service
[Unit]
Description=replicated

[Service]
ExecStart=/bin/sleep 1d
