# shipyard/shipyard/py/garage/examples/echod/files/echod.service
[Unit]
Description=echod

[Service]
Slice=machine.slice
ExecStart=/usr/bin/rkt run --net=host $IMAGE_echod -- -m echod
KillMode=mixed
Restart=always
